# vrchat-status-push.service
# 此文件为模板，推荐使用 install.sh 自动安装（自动检测路径和 Python 位置）
# 手动安装时需将路径替换为实际值
[Unit]
Description=VRChat Status Webhook Push
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/vrchat-status-push/main.py -c /opt/vrchat-status-push/config.json -s /opt/vrchat-status-push/data/state.json
WorkingDirectory=/opt/vrchat-status-push
Restart=always
RestartSec=30
User=nobody  # 可按需改为你的用户名

[Install]
WantedBy=multi-user.target
